# DuckDNS.timer
[Unit]
Description=Timing DuckDNS update

[Timer]
OnCalendar=*-*-* *:5:00
OnStartupSec=2m

[Install]
WantedBy=timers.target
